# docker/crontab
# Crontab для AutoRewrite
# Запуск парсинга контента каждый час
0 * * * * php /var/www/html/cron/parse.php >> /var/log/cron.log 2>&1

# Запуск реврайта контента каждые 2 часа
0 */2 * * * php /var/www/html/cron/rewrite.php >> /var/log/cron.log 2>&1

# Запуск публикации контента каждые 3 часа
0 */3 * * * php /var/www/html/cron/post.php >> /var/log/cron.log 2>&1

# Пустая строка в конце файла
